# docker/supervisor.conf
# docker/supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log ; Jalur log untuk supervisor, bisa juga diarahkan ke stdout/stderr
pidfile=/var/run/supervisord.pid
loglevel=info ; Bisa diubah ke debug untuk informasi lebih detail saat troubleshooting

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F 
; Opsi '-F' (foreground) penting agar Supervisor bisa mengelolanya dengan benar.
; Jangan gunakan versi spesifik seperti php-fpm8.4 kecuali Anda 100% yakin itu nama executable di base image Anda.
; Untuk image php:8.2-fpm-alpine, /usr/local/sbin/php-fpm adalah yang benar.
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=www-data ; Jalankan sebagai user www-data untuk konsistensi izin

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
; Opsi '-g "daemon off;"' penting agar Nginx berjalan di foreground.
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
